# server/config/ekg3000server.conf
description "EKG3000 Server Service"
author  "Noah Huetter noahhuetter@gmail.com"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/ekg3000/pro3-web/server
exec /home/ekg3000/pro3-web/server/server.py
respawn

